# build/djcelery
### BEGIN INIT INFO
# Provides: djcelery
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start django celery daemon at boot time
# Description: Enable service provided by daemon.
### END INIT INFO
 
# chkconfig: 345 20 80
# description: django-celery daemon
# processname: djcelery

# /etc/init.d/djcelery

case "$1" in
  start)
        echo -e "Starting djcelery service\n"
        #To run it as root:
        python /home/meuprojeto/wsgi/dev/meuprojeto/async_tasks/start_djcelery.py
	;;
  stop)
        echo -e "Stopping djcelery\n"
        #To run it as root:
        #/path/to/command/to/stop/counter
        #Or to run it as some other user:
        #/bin/su - username -c /path/to/command/to/stop/counter
        #echo "."
        kill -9 `ps -ef | grep -i 'meuprojeto/manage.py celeryd' | grep -v grep | awk '{print $2}'`
	;;

  *)
        echo "Usage: /sbin/service counter {start|stop}"
        exit 1
esac

exit 0


